# edge/os-baseline/systemd/atlasmesh-tpm-manager.service
[Unit]
Description=AtlasMesh TPM Manager
Documentation=https://docs.atlasmesh.com/tpm-manager
After=tpm2-abrmd.service
Wants=tpm2-abrmd.service

[Service]
Type=simple
ExecStart=/usr/local/bin/atlasmesh-tpm-manager
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=15

# Security settings
User=atlasmesh
Group=atlasmesh
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=false
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true

# TPM device access
DeviceAllow=/dev/tpm0 rw
DeviceAllow=/dev/tpmrm0 rw
DevicePolicy=closed

# File system access
ReadWritePaths=/var/lib/atlasmesh/certs
ReadOnlyPaths=/etc/atlasmesh

# Environment
Environment=ATLASMESH_HOME=/var/lib/atlasmesh
Environment=ATLASMESH_CONFIG=/etc/atlasmesh
Environment=TPM_DEVICE=/dev/tpm0
EnvironmentFile=-/etc/atlasmesh/tpm-manager.env

# Resource limits
LimitNOFILE=1024
LimitNPROC=512
MemoryMax=256M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
